# pi/pi-access-point/etc/rc.local
iptables-restore < /etc/iptables.ipv4.nat
ip -br addr
exit 0
